# debian/DEBIAN/prerm
#!/bin/bash
systemctl stop ikea-ad.service
systemctl disable ikea-ad.service
